# samples/docker/renderer/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
stdout_events_enabled=true
stderr_events_enabled=true

[program:syslog-ng]
command=/usr/sbin/syslog-ng --foreground --no-caps

[program:static-renderer]
command=/usr/bin/fastcgistart2.sh yandex-maps-staticrenderer
